# local/docker/php-fpm/supervisor/supervisord.conf
[supervisord]
nodaemon=true
[program:cron]
command=cron -f -L 15
[program:php]
command=php-fpm --allow-to-run-as-root